# infra/terraform/templates/lightspeed-proxy.service
[Unit]
Description=LightSpeed Proxy — UDP game latency optimizer
Documentation=https://github.com/ShibbityShwab/lightspeed
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30

# Native binary — no Docker overhead (~568KB vs ~200MB for Docker daemon)
ExecStart=/usr/local/bin/lightspeed-proxy \
  --config /etc/lightspeed/proxy.toml \
  --data-bind 0.0.0.0:4434 \
  --control-bind 0.0.0.0:4433 \
  --health-bind 0.0.0.0:8080

# Run as dedicated unprivileged user (systemd dynamic allocation)
DynamicUser=true

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/etc/lightspeed
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
